# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:api]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/api
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/api.err.log
stdout_logfile=/var/log/supervisor/api.out.log
environment=PYTHONUNBUFFERED=1

[program:ui]
command=serve -s /app/ui/dist -l 3000
directory=/app/ui
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/ui.err.log
stdout_logfile=/var/log/supervisor/ui.out.log
